# root/etc/s6-overlay/s6-rc.d/init-fireqos/up
#!/command/execlineb -P

with-contenv
pipeline {
    fdmove -c 2 1
    /usr/lib/mihomo/init-fireqos.sh
}
mawk -W interactive "{printf \"%c[34m[init-fireqos]%c[0m %s\\n\", 27, 27, $0}"
