# root/etc/s6-overlay/s6-rc.d/init-network/up
#!/command/execlineb -P

with-contenv
pipeline {
    fdmove -c 2 1
    /usr/lib/mihomo/init-network.sh
}
mawk -W interactive "{printf \"%c[32m[init-network]%c[0m %s\\n\", 27, 27, $0}"
